# applications/luci-olsr/root/etc/uci-defaults/luci-olsr
#!/bin/sh

uci batch <<-EOF
	add ucitrack olsrd
	set ucitrack.@olsrd[-1].init=olsrd
	commit ucitrack
EOF
