# themes/freifunk-hannover/root/etc/uci-defaults/luci-theme-freifunk-hannover
uci batch <<-EOF
        set luci.themes.Freifunk_Hannover=/luci-static/freifunk-hannover
        commit luci
EOF
