# modules/admin-full/root/etc/init.d/luci_bwc
#!/bin/sh /etc/rc.common

START=95
STOP=95

BWC=/usr/bin/luci-bwc

start() {
	$BWC -d
}

stop() {
	killall ${BWC##*/}
}
