# themes/luci-theme-purple/root/etc/uci-defaults/30_luci-theme-purple
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Purple=/luci-static/purple
	set luci.main.mediaurlbase=/luci-static/purple
	commit luci
EOF
exit 0
