# applications/luci-p910nd/root/etc/uci-defaults/luci-p910nd
#!/bin/sh

uci batch <<-EOF
	add ucitrack p910nd
	set ucitrack.@p910nd[-1].init=p910nd
	commit ucitrack
EOF
