# i18n/romanian/root/etc/uci-defaults/luci-i18n-romanian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.ro=Română
	commit luci
EOF

