# applications/luci-app-adblock-fast/root/etc/uci-defaults/40_luci-adblock-fast
#!/bin/sh
rm -rf /var/luci-modulecache/; rm -f /var/luci-indexcache;
[ -x /etc/init.d/rpcd ] && /etc/init.d/rpcd reload;
exit 0
