# themes/fledermaus/root/etc/uci-defaults/luci-theme-fledermaus
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Fledermaus=/luci-static/fledermaus
        commit luci
EOF
	
