# applications/luci-app-advanced-reboot/root/etc/uci-defaults/40_luci-advanced-reboot
#!/bin/sh

/etc/init.d/rpcd reload
rm -rf /var/luci-modulecache/; rm -f /var/luci-indexcache;
exit 0
