# themes/luci-theme-nightstrap/root/etc/uci-defaults/30_luci-theme-nightstrap
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap=/luci-static/nightstrap
	set luci.main.mediaurlbase=/luci-static/nightstrap
	commit luci
EOF
exit 0
