# themes/luci-theme-freifunk-bno/root/etc/uci-defaults/luci-theme-freifunk-bno
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Freifunk_BNO=/luci-static/freifunk-bno
	set luci.main.mediaurlbase=/luci-static/freifunk-bno
        commit luci
EOF
