# i18n/catalan/root/etc/uci-defaults/luci-i18n-catalan
#!/bin/sh
uci batch <<-EOF
	set luci.languages.ca=Catalan
	commit luci
EOF
	
