# i18n/japanese/root/etc/uci-defaults/luci-i18n-japanese
#!/bin/sh
uci batch <<-EOF
	set luci.languages.ja=Japanese
        commit luci
EOF
	
