# contrib/package/freifunk-policyrouting/files/etc/hotplug.d/iface/30-policyrouting
[ "$INTERFACE" != "wan" ] && exit 0

case $ACTION in
        ifup)
		pr=`uci get freifunk-policyrouting.pr.enable`
		# check if ipv6 is enabled:
		[ -f /proc/net/ipv6_route ] && has_ipv6=1

		if [ $pr = "1" ]; then
			logger -t policyrouting "Starting policy routing on $INTERFACE"

			# Setup new tables
			tables="/etc/iproute2/rt_tables"
			if [ -z "`grep "111" $tables`" ]; then
				echo "111 olsr" >> $tables
			fi
			if [ -z "`grep "112" $tables`" ]; then
				echo "112 olsr-default" >> $tables
			fi

			# Make sure Rt_tables in olsrd are in place
			if [ ! "`uci -q get olsrd.@olsrd[0].RtTable`" == "111" ] || [ ! "`uci -q get olsrd.@olsrd[0].RtTableDefault`" == "112" ]; then
				uci set olsrd.@olsrd[0].RtTable='111'
				uci set olsrd.@olsrd[0].RtTableDefault='112'
				uci commit
				/etc/init.d/olsrd restart
			fi

			# Disable dyn_gw and dyngw_plain
			dyngwlib=`uci show olsrd |grep dyn_gw.so |awk {' FS="."; print $1"."$2 '}`
			if [ -n "$dyngwlib" ]; then
				uci set $dyngwlib.ignore=1
				uci commit
			fi

			dyngwplainlib=`uci show olsrd |grep dyn_gw_plain |awk {' FS="."; print $1"."$2 '}`
			if [ -n "$dyngwplainlib" ]; then
				uci set $dyngwplainlib.ignore=1
				uci commit
			fi
	
			gw=`uci -p /var/state get network.wan.gateway`
			netmask=`uci -p /var/state get network.wan.netmask`
			if [ -z "$netmask" ]; then
				NETMASK="255.255.255.255"
			fi

			if  [ -n "`uci -p /var/state get network.wan.ifname`" ]; then
				device=`uci -p /var/state get network.wan.ifname`
			else
				device=`uci -p /var/state get network.wan.device`
			fi

			eval `ipcalc.sh $gw $netmask`

			test -n "`ip r s t default`" && ip r d default t default
			test -n "`ip r s |grep default`" && ip route del default
			ip route add $NETWORK/$NETMASK dev $device table default
			ip route add default via $gw dev $device table default

			#if [ "$has_ipv6" = 1 ]; then
			#	local ip6gw=$(ip -6 r |grep default |cut -d " " -f 3)
			#	test -n "`ip -6 r s t default`" && ip -6 r d default t default
			#	test -n "`ip -6 r s |grep default`" && ip -6 route del default
			#	ip -6 r a $ip6gw via $ip6gw dev $dev table default
			#	ip -6 route add default via $ip6gw dev $device table default
			#fi


			ip rule del lookup main
			ip rule add fwmark 1 lookup olsr-default
			ip rule add lookup main
			ip rule add lookup olsr
			if [ "$has_ipv6" = 1 ]; then
	                        ip -6 rule del lookup main
	                        ip -6 rule add fwmark 1 lookup olsr-default prio 16385
	                        ip -6 rule add lookup main prio 16383
	                        ip -6 rule add lookup olsr prio 16380
			fi
		else
			# Remove custom routing tables from olsrd
			if [ "`uci -q get olsrd.@olsrd[0].RtTable`" == "111" ] || [ "`uci -q get olsrd.@olsrd[0].RtTableDefault`" == "112" ]; then
				uci delete olsrd.@olsrd[0].RtTable
				uci delete olsrd.@olsrd[0].RtTableDefault
				uci commit
				/etc/init.d/olsrd restart
			fi
		fi
	;;

	ifdown)
	logger -t policyrouting "Deleting policy rules for $INTERFACE"
	ip rule del fwmark 1 lookup olsr-default > /dev/null 2>&1
	ip rule del lookup olsr	> /dev/null 2>&1
	if [ "$has_ipv6" = 1 ]; then
		ip -6 rule del fwmark 1 lookup olsr-default > /dev/null 2>&1
		ip -6 rule del lookup olsr	> /dev/null 2>&1
	fi
	;;
esac
