# applications/luci-app-ushare/root/etc/uci-defaults/40_luci-ushare
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ushare[-1]
	add ucitrack ushare
	set ucitrack.@ushare[-1].init=ushare
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
