# contrib/package/freifunk-gwcheck/root/etc/uci-defaults/freifunk-gwcheck
#!/bin/sh
tables="/etc/iproute2/rt_tables"
[ -z "`grep "gw-check" $tables`" ] && echo "200 gw-check" >> $tables
[ -z "`grep "ff_olsr_test_gw.sh" /etc/crontabs/root`" ] && echo "* * * * *	/usr/sbin/ff_olsr_test_gw.sh" >> /etc/crontabs/root
/etc/init.d/cron restart
