# contrib/package/luci-splash/src/luci_splash.cron
*/15 * * * *	/usr/lib/luci-splash/sync.lua