# themes/openwrt-light/root/etc/uci-defaults/luci-theme-openwrtlight
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenWRT_Light=/luci-static/openwrt-light
        commit luci
EOF
	
