# themes/openwrt.org-oxygen/root/etc/uci-defaults/luci-theme-openwrt-oxygen
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenWrtOxygen=/luci-static/openwrt.org-oxygen
        commit luci
EOF
	
