# themes/luci-theme-openwrtcn/root/etc/uci-defaults/luci-theme-openwrtcn
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenWrtCN=/luci-static/openwrtcn
	commit luci
EOF
	
