# i18n/norwegian/etc/uci-defaults/luci-i18n-norwegian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.no=Norwegian
	commit luci
EOF
