# i18n/italian/root/etc/uci-defaults/luci-i18n-italian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.it=Italiano
        commit luci
EOF

