# themes/freifunk/root/etc/uci-defaults/luci-theme-freifunk
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Freifunk=/luci-static/freifunk
        commit luci
EOF
	
