# themes/luci-theme-oxygen/root/etc/uci-defaults/luci-theme-oxygen
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Oxygen=/luci-static/oxygen
        commit luci
EOF
	
