# themes/luci-theme-freifunk-generic/root/etc/uci-defaults/30_luci-theme-freifunk-generic
uci batch <<-EOF
	set luci.themes.Freifunk_Generic=/luci-static/freifunk-generic
	set luci.main.mediaurlbase=/luci-static/freifunk-generic
        commit luci
EOF
