# themes/luci-theme-freifunk-generic/root/etc/uci-defaults/luci-theme-freifunk-generic
uci batch <<-EOF
        set luci.themes.Freifunk_Generic=/luci-static/freifunk-generic
        commit luci
EOF
