# applications/luci-app-libreswan/root/etc/uci-defaults/802-luci-app-libreswan
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@libreswan[-1]
	add ucitrack libreswan
	set ucitrack.@libreswan[-1].init=ipsec
	commit ucitrack
EOF
