# i18n/russian/root/etc/uci-defaults/luci-i18n-russian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.ru=Русский
        commit luci
EOF

