# contrib/package/ffluci-splash/src/luci_splash.cron
#!/bin/sh
[ "$(date +%M | cut -c2)" == "5" ] && luci-splash sync