# themes/luci-theme-openwrt/root/etc/uci-defaults/30_luci-theme-openwrt
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenWrt=/luci-static/openwrt.org
        commit luci
EOF
