# applications/luci-statistics/root/etc/uci-defaults/luci-statistics
#!/bin/sh

uci set luci.uci_oncommit.luci_statistics='/etc/init.d/luci_statistics restart'
uci commit luci
