# i18n/chinese/root/etc/uci-defaults/luci-i18n-chinese
#!/bin/sh
uci batch <<-EOF
	set luci.languages.zh_cn=chinese
	commit luci
EOF
