# i18n/french/root/etc/uci-defaults/luci-i18n-french
#!/bin/sh
uci batch <<-EOF
	set luci.languages.fr='Français'
        commit luci
EOF

