# applications/luci-asterisk/root/etc/uci-defaults/luci-asterisk
#!/bin/sh

uci batch <<-EOF
	add ucitrack asterisk
	set ucitrack.@olsrd[-1].init=asterisk
	commit ucitrack
EOF
