# themes/luci-theme-xeye/root/etc/uci-defaults/luci-theme-lucixeye
#!/bin/sh
uci batch <<-EOF
	set luci.themes.LuciXEyE=/luci-static/LuciXEyE
        commit luci
EOF
	
