# applications/luci-app-babeld/root/etc/uci-defaults/40_luci-babeld
#!/bin/sh

uci set babeld.@general[0].ubus_bindings='true'

return 0
