# contrib/package/freifunk-common/files/etc/uci-defaults/freifunk
uci set uhttpd.main.rfc1918_filter=0
uci commit uhttpd
