# applications/freifunk-community/root/etc/uci-defaults/luci-freifunk-community
#!/bin/sh
uci batch <<-EOF
        set luci.main.mediaurlbase=/luci-static/fledermaus
        set luci.main.lang=de
        commit luci
EOF
	
