# modules/freifunk/root/etc/rc.local
# /etc/rc.local
# Only executed on boot, put your custom startup commands here
