# i18n/hungarian/root/etc/uci-defaults/luci-i18n-hungarian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.hu=Magyar
	commit luci
EOF
