# applications/luci-ushare/root/etc/uci-defaults/luci-ushare
#!/bin/sh

uci batch <<-EOF
	add ucitrack ushare
	set ucitrack.@ushare[-1].init=ushare
	commit ucitrack
EOF
