# i18n/malay/root/etc/uci-defaults/luci-i18n-malay
#!/bin/sh
uci batch <<-EOF
	set luci.languages.ms='Bahasa Melayu'
        commit luci
EOF

