# i18n/portuguese-brazilian/root/etc/uci-defaults/luci-i18n-portuguese_brazilian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.pt_br="Português (Brasil)"
        commit luci
EOF
