# applications/luci-app-shadowsocks-libev/root/etc/uci-defaults/40_luci-shadowsocks-libev
#!/bin/sh

rm -f /tmp/luci-indexcache
mkdir -p /etc/shadowsocks-libev
/etc/init.d/rpcd reload
exit 0
