# i18n/vietnamese/root/etc/uci-defaults/luci-i18n-vietnamese
#!/bin/sh
uci batch <<-EOF
	set luci.languages.vi=Vietnamese
	commit luci
EOF
	
