# i18n/vietnamese/root/etc/uci-defaults/luci-i18n-vietnamese
#!/bin/sh
uci batch <<-EOF
	set luci.languages.vi="Tiếng Việt"
	commit luci
EOF

