# applications/luci-app-syncdial/root/etc/hotplug.d/iface/01-dialcheck
#!/bin/sh
[ "$ACTION" = "ifdown" ] && pppconnectcheck &


