# i18n/norwegian/root/etc/uci-defaults/luci-i18n-norwegian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.no=Norsk
        commit luci
EOF

