# i18n/spanish/root/etc/uci-defaults/luci-i18n-spanish
#!/bin/sh
uci batch <<-EOF
	set luci.languages.es=Español
        commit luci
EOF

