# applications/luci-radvd/root/etc/uci-defaults/luci-radvd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@radvd[-1]
	add ucitrack radvd
	set ucitrack.@radvd[-1].init=radvd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
