# applications/luci-p2pblock/root/etc/uci-defaults/luci-p2pblock
#!/bin/sh

uci batch <<-EOF 2>/dev/null
	add ucitrack freifunk_p2pblock
	set ucitrack.@freifunk_p2pblock[-1].init=freifunk-p2pblock
	commit ucitrack
EOF

exit 0
