# applications/luci-app-arpbind/root/etc/uci-defaults/luci-arpbind
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@arpbind[-1]
	commit ucitrack
EOF

exit 0
