# applications/luci-app-https-dns-proxy/root/etc/uci-defaults/40_luci-https-dns-proxy
#!/bin/sh
rm -rf /var/luci-modulecache/
exit 0
