# i18n/german/root/etc/uci-defaults/luci-i18n-german
#!/bin/sh
uci batch <<-EOF
	set luci.languages.de=Deutsch
        commit luci
EOF

