# i18n/ukrainian/root/etc/uci-defaults/luci-i18n-ukrainian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.en=English
        commit luci
EOF

