# i18n/ukrainian/root/etc/uci-defaults/luci-i18n-ukrainian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.uk='украї́нська'
        commit luci
EOF

