# applications/luci-app-adblock/root/etc/uci-defaults/95-luci-app-adblock-housekeeping
#!/bin/sh
rm -f /var/luci-indexcache.*.json
[ -x "/etc/init.d/rpcd" ] && /etc/init.d/rpcd reload
exit 0
