# i18n/greek/root/etc/uci-defaults/luci-i18n-greek
#!/bin/sh
uci batch <<-EOF
	set luci.languages.el=Ελληνικά
        commit luci
EOF

