# applications/luci-app-irqbalance/root/etc/uci-defaults/90_luci-irqbalance
#!/bin/sh

uci -q batch <<-EOF >/dev/null
    delete ucitrack.@irqbalance[-1]
    add ucitrack irqbalance
    set ucitrack.@irqbalance[-1].init=irqbalance
    commit ucitrack
EOF

exit 0
