# i18n/hebrew/root/etc/uci-defaults/luci-i18n-hebrew
#!/bin/sh
uci batch <<-EOF
	set luci.languages.he=עִבְרִית
	commit luci
EOF

