# applications/luci-hd_idle/root/etc/uci-defaults/luci-hd_idle
#!/bin/sh

uci batch <<-EOF
       add ucitrack hd-idle
       set ucitrack.@hd-idle[-1].init=hd-idle
       commit ucitrack
EOF
