# themes/openwrt.org/root/etc/uci-defaults/luci-theme-openwrt
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenWRT=/luci-static/openwrt.org
        commit luci
EOF
	
