# applications/luci-app-zerotier/root/etc/uci-defaults/40_luci-zerotier
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@zerotier[-1]
	commit ucitrack

	delete firewall.zerotier
	commit firewall
EOF

rm -f /tmp/luci-indexcache
exit 0
