# contrib/package/luci-addons/dist/etc/crontabs/root
* * * * *  /usr/bin/run-parts /etc/cron.minutely
17 * * * *       /usr/bin/run-parts /etc/cron.hourly
25 6 * * *       /usr/bin/run-parts /etc/cron.daily
