# applications/luci-app-acme/root/etc/uci-defaults/60_luci-acme
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@acme[-1]
	add ucitrack acme
        add_list ucitrack.@bcp38[0].init=acme
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
