# applications/luci-uvc_streamer/root/etc/uci-defaults/luci-uvc_streamer
#!/bin/sh

uci set luci.uci_oncommit.uvc_streamer='/sbin/luci-reload uvc-streamer'
uci commit luci
