# applications/luci-uvc_streamer/root/etc/uci-defaults/luci-uvc_streamer
#!/bin/sh

uci batch <<-EOF
	add ucitrack uvc-streamer
	set ucitrack.@uvc-streamer[-1].init=uvc-streamer
	commit ucitrack
EOF
