# applications/luci-meshwizard/root/etc/uci-defaults/meshwizard
#!/bin/sh
uci batch <<EOF
set ucitrack.meshwizard="meshwizard"
set ucitrack.meshwizard.exec="/etc/init.d/wizard boot"
commit ucitrack
EOF

rm -f /tmp/luci-indexcache
