# applications/luci-upnp/root/etc/uci-defaults/luci-upnp
#!/bin/sh
/etc/init.d/miniupnpd stop
/etc/init.d/miniupnpd disable
