# contrib/package/ffluci-system-addons/src/root.crontab
0-59/1 * * * *  /usr/bin/run-parts /etc/cron.minutely
0 * * * *       /usr/bin/run-parts /etc/cron.hourly
0 0 * * *       /usr/bin/run-parts /etc/cron.daily
