# applications/luci-app-uhttpd/root/etc/uci-defaults/luci-uhttpd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@uhttpd[-1]
	add ucitrack uhttpd
	set ucitrack.@uhttpd[-1].init=uhttpd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
