# applications/luci-polipo/root/etc/uci-defaults/luci-polipo
#!/bin/sh

uci batch <<-EOF
	add ucitrack polipo
	set ucitrack.@polipo[-1].init=polipo
	commit ucitrack
EOF
