# i18n/swedish/root/etc/uci-defaults/luci-i18n-swedish
#!/bin/sh
uci batch <<-EOF
	set luci.languages.sv=Svenska
        commit luci
EOF

